# camlistored.conf
description	"Start camlistore"

start on runlevel [2345]

script
    export USER=root
    exec /usr/bin/camlistored
end script
